# agent/testdata/test.service
[Unit]
Description=test
[Service]
ExecStart=/usr/bin/sleep inf
[Install]
WantedBy=default.target