# testdata/test.service
[Unit]
Description=%p

[Service]
ExecStart=/usr/bin/sleep inf

[Install]
WantedBy=multi-user.target